# roles/proxy/templates/node-discovery.service.j2
[Unit]
Description=Sovereign node-discovery (PostgreSQL -> OpenResty backend cache)
After=network-online.target openresty.service
Wants=network-online.target
Requires=openresty.service

[Service]
Type=simple
EnvironmentFile=/etc/node-discovery.env
ExecStart=/usr/local/bin/node-discovery \
  --database-url ${DATABASE_URL} \
  --output-file ${OUTPUT_FILE} \
  --max-age-millis {{ proxy_node_discovery_max_age_millis }} \
  --nginx-binary /usr/local/openresty/nginx/sbin/nginx \
  --metrics-port {{ proxy_node_discovery_metrics_port }}
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
